# deploy/etc_init/campjs-snake.conf
# upstart service definition

description "run campjs-snake server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/nodejs /opt/campjs-snake/index.js

